# pi/autostart/calendar-kiosk.desktop
[Desktop Entry]
Type=Application
Name=Calendar Kiosk
Exec=chromium \
  --kiosk \
  --noerrdialogs \
  --disable-infobars \
  --disable-session-crashed-bubble \
  --disable-restore-session-state \
  --no-first-run \
  --disable-sync \
  --user-data-dir=/tmp/chromium-kiosk \
  --check-for-update-interval=31536000 \
  http://localhost/
X-GNOME-Autostart-enabled=true
